# etc/supervisord.d/opentsdb.ini
[program:opentsdb]
command = /opt/qnib/opentsdb/bin/start_opentsdb.sh
redirect_stderr=true
stdout_logfile=/var/log/supervisor/%(program_name)s.log
